# .config/systemd/user/mullvad-vpn.service
[Unit]
Description=Mullvad VPN GUI

[Service]
ExecStart=/usr/bin/mullvad-vpn
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=default.target
